# deploy/systemd/qltracker-provisioner.service
[Unit]
Description=QLTracker Pickup Provisioner
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=qltracker
Group=qltracker
WorkingDirectory=/opt/qltracker-provisioner/app
EnvironmentFile=/opt/qltracker-provisioner/app/.env
ExecStart=/usr/bin/node /opt/qltracker-provisioner/app/dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
