# apps/provisioner/deploy/systemd/qltracker-provisioner.service
[Unit]
Description=QLTracker Pickup Provisioner
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=qltracker
Group=qltracker
WorkingDirectory=/opt/qltracker/apps/provisioner
EnvironmentFile=/opt/qltracker/apps/provisioner/.env
ExecStart=/usr/bin/env node /opt/qltracker/apps/provisioner/dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
